# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:appium]
user=jenkins
command=/bin/bash -c "appium"
